# cron.d/mrtg
*/5 * * * *		/bin/bash /opt/mrtg/cron5m.sh
