# cron.d/mrtg
*/5 * * * *		/bin/bash /etc/mrtg/cron5m.sh
